# buildroot-external/package/bluetooth-bcm43xx/bthelper@.service
[Unit]
Description=Raspberry Pi bluetooth helper
Requires=bluetooth.service
After=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/bthelper %I
