# buildroot-external/rootfs-overlay/usr/lib/systemd/system/systemd-time-wait-sync.service.d/network-online.conf
[Unit]
After=network-online.target
Wants=network-online.target

[Install]
WantedBy=
WantedBy=time-sync.target
