# buildroot-external/board/jethome/rootfs-overlay/usr/lib/udev/rules.d/mac-set.rules
SUBSYSTEM=="net",KERNEL=="eth0",RUN+="/usr/lib/jethome/jsetmac %k"
SUBSYSTEM=="net",KERNEL=="wlan0",RUN+="/usr/lib/jethome/jsetmac %k"