# buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassos-zram.target
[Unit]
Description=HassOS zram targets
Documentation=man:systemd.target(5)

[Install]
WantedBy=multi-user.target
