# buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassos-bind.target
[Unit]
Description=HassOS overlay targets
Documentation=man:systemd.target(5)
Before=sysinit.target

[Install]
WantedBy=multi-user.target
