# buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassos-hardware.target
[Unit]
Description=HassOS hardware targets
Documentation=man:systemd.target(5)
