# buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassio-expand.service
[Unit]
Description=Hassio data resizing
DefaultDependencies=no
Before=mnt-data.mount

[Service]
Type=oneshot
ExecStart=/sbin/hassio-expand
RemainAfterExit=true

[Install]
WantedBy=local-fs.target
