# buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassio-zram.target
[Unit]
Description=Hassio zram targets
Documentation=man:systemd.target(5)
