# buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassos-apparmor.service
[Unit]
Description=HassOS AppArmor
Wants=hassos-supervisor.service
Before=docker.service hassos-supervisor.service
RequiresMountsFor=/mnt/data

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/libexec/hassos-apparmor

[Install]
WantedBy=multi-user.target
