# buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassos-bind.target
[Unit]
Description=HassOS overlay targets
Documentation=man:systemd.target(5)
After=local-fs.target
Before=sysinit.target
