# buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassio-bind.target
[Unit]
Description=Hassio overlay targets
Documentation=man:systemd.target(5)
