# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mnt-config.mount
[Unit]
Description=HassOS config partition

[Mount]
What=/dev/disk/by-label/CONFIG
Where=/mnt/config
Type=auto
Options=ro
