# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mnt-boot.mount
[Unit]
Description=HassOS boot partition
DefaultDependencies=no
Before=umount.target local-fs.target
After=systemd-fsck@dev-disk-by\x2dlabel-hassos\x2dboot.service
Wants=systemd-fsck@dev-disk-by\x2dlabel-hassos\x2dboot.service
Conflicts=umount.target

[Mount]
What=/dev/disk/by-label/hassos-boot
Where=/mnt/boot
Type=auto
Options=sync

[Install]
WantedBy=local-fs.target
