# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mnt-config.mount
[Unit]
Description=HassOS config partition

[Mount]
What=LABEL=config
Where=/mnt/config
Type=auto
Options=ro
