# buildroot-external/rootfs-overlay/usr/lib/systemd/system/rauc-good.timer
[Unit]
Description=Rauc mark boot partition as good

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
