# buildroot-external/rootfs-overlay/usr/lib/systemd/system/raucdb-update.service
[Unit]
Description=Ensure rauc.db contains version information
After=rauc.service
Before=hassos-supervisor.service
RequiresMountsFor=/mnt/data

[Service]
Type=oneshot
ExecStart=/usr/libexec/raucdb-update

[Install]
WantedBy=multi-user.target
