# buildroot-external/rootfs-overlay/usr/lib/systemd/system/rauc-good.timer
[Unit]
Description=HassOS mark boot partition as good

[Timer]
OnBootSec=30sec

[Install]
WantedBy=timers.target
