# buildroot-external/rootfs-overlay/usr/lib/systemd/system/rauc-mark.timer
[Unit]
Description=HassOS mark boot partition as good

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
