# buildroot-external/rootfs-overlay/usr/lib/systemd/system/rauc-bad.service
[Unit]
Description=Rauc mark bad
Requires=rauc.service
After=rauc.service
RefuseManualStart=true
RefuseManualStop=true

[Service]
Type=oneshot
ExecStart=/usr/bin/rauc status mark-bad
ExecStartPost=/usr/bin/systemctl reboot

[Install]
WantedBy=multi-user.target
