# buildroot-external/rootfs-overlay/usr/lib/systemd/system/zram-var.service
[Unit]
Description=HassOS ZRAM var
DefaultDependencies=no
Before=var.mount

[Service]
Type=oneshot
ExecStart=/usr/sbin/hassos-zram -t fs -m var -s 32M
RemainAfterExit=true
StandardOutput=null

[Install]
WantedBy=hassos-zram.target
