# buildroot-external/overlay/etc/systemd/system/getty@tty1.service.d/mingetty.conf
[Service]
ExecStart=
ExecStart=-/sbin/mingetty --noclear %I
