# buildroot-external/rootfs-overlay/usr/lib/systemd/system/serial-getty@.service.d/hassos.conf
[Unit]
Wants=hassos-supervisor.service
After=hassos-supervisor.service
