# buildroot-external/package/xe-guest-utilities/proc-xen.mount
[Unit]
Description=Mount /proc/xen files
ConditionPathExists=/proc/xen
RefuseManualStop=true

[Mount]
What=xenfs
Where=/proc/xen
Type=xenfs
