# buildroot/package/irqbalance/irqbalance.service
[Unit]
Description=irqbalance daemon
After=syslog.target

[Service]
RuntimeDirectory=irqbalance
ExecStart=/usr/sbin/irqbalance --foreground

[Install]
WantedBy=multi-user.target
