# buildroot-external/rootfs-overlay/usr/lib/systemd/system/qemu-guest.service.d/haos.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/qemu-ga -m virtio-serial -p /dev/virtio-ports/org.qemu.guest_agent.0 --fsfreeze-hook=/usr/libexec/haos-freeze-hook
