# buildroot-external/rootfs-overlay/etc/systemd/system/dropbear.service.d/hassos.conf
[Unit]
RequiresMountsFor=/etc/dropbear
ConditionFileNotEmpty=/root/.ssh/known_hosts

[Service]
ExecStartPre=
ExecStart=
ExecStart=/usr/sbin/dropbear -F -R -E -p 22222 -s
