# buildroot-external/rootfs-overlay/etc/systemd/system/dropbear.service.d/hassio.conf
[unit]
Requires=etc-dropbear.mount
After=etc-dropbear.mount

[Service]
ExecStartPre=
ExecStart=-/usr/sbin/dropbear -F -R -E -p 22222 -s
