# buildroot-external/rootfs-overlay/etc/systemd/system/dropbear.service.d/nm.conf
[Unit]
Wants=network-online.target
After=network-online.target
