# buildroot-external/rootfs-overlay/usr/lib/systemd/system/systemd-localed.service.d/vconsole.conf
[Service]
Environment=SYSTEMD_ETC_VCONSOLE_CONF=/mnt/overlay/etc/vconsole.conf
ReadWritePaths=/etc /mnt/overlay/etc
