# buildroot-external/rootfs-overlay/etc/systemd/system/NetworkManager.service.d/hassos.conf
[Unit]
RequiresMountsFor=/etc/hostname /etc/hosts /etc/NetworkManager/system-connections
