# buildroot-external/rootfs-overlay/etc/systemd/system/dbus-broker.service.d/haos.conf
[Service]
PrivateTmp=false
PrivateDevices=false
