# buildroot-external/rootfs-overlay/etc/systemd/system/systemd-timesyncd.service.d/hassos.conf
[Unit]
RequiresMountsFor=/var/lib/systemd
After=network-online.target

[Install]
WantedBy=
WantedBy=time-sync.target
