# buildroot-external/rootfs-overlay/etc/systemd/system/systemd-timesyncd.service.d/ro.conf
[Service]
ExecStartPre=-/bin/sh -c '[ ! -d /var/lib/systemd/timesync ] && rm /var/lib/systemd/timesync'
DynamicUser=no
