# buildroot-external/rootfs-overlay/etc/systemd/system/apparmor.service.d/hassio.conf
[Unit]
RequiresMountsFor=/etc/apparmor.d/containers
