# buildroot-external/rootfs-overlay/etc/systemd/system/systemd-time-wait-sync.service.d/network-online.conf
[Unit]
Wants=network-online.target

[Install]
WantedBy=time-sync.target
