# buildroot-external/rootfs-overlay/etc/systemd/system/systemd-modules-load.service.d/hassos.conf
[Unit]
RequiresMountsFor=/etc/modules-load.d /etc/modprobe.d
