# buildroot-external/rootfs-overlay/etc/systemd/system/getty@tty1.service.d/mingetty.conf
[Service]
ExecStart=-/sbin/mingetty --noclear %I
