# buildroot-external/overlay/etc/systemd/system/mnt-boot.mount
[Unit]
Description=Hassio boot partition

[Mount]
What=LABEL=hassio-boot
Where=/mnt/boot
Type=auto

[Install]
WantedBy=multi-user.target
