# buildroot-external/rootfs-overlay/etc/systemd/system/serial-getty@.service.d/docker.conf
[Unit]
Requires=docker.service
After=docker.service
