# buildroot-external/board/jethome/jethub-d1/rootfs-overlay/usr/lib/systemd/system/zram-swap.service
[Unit]
Description=HassOS File swap for JetHub D1
DefaultDependencies=no
Requires=mnt-data.mount
After=mnt-data.mount hassos-expand.service
Wants=mnt-data.mount hassos-expand.service
#RefuseManualStart=true
#RefuseManualStop=true

[Service]
Type=oneshot
<<<<<<< HEAD
ExecStart=/usr/lib/jethome/swap_file_maker
TimeoutStartSec=15min
RemainAfterExit=true
#StandardOutput=null

[Install]
WantedBy=hassos-zram.target
