# buildroot-external/rootfs-overlay/usr/lib/systemd/system/docker.service.d/failure.conf
[Unit]
FailureAction=reboot

[Service]
ExecStartPre=/usr/libexec/docker-disk-check
