# buildroot-external/rootfs-overlay/usr/lib/systemd/system/docker.service.d/hassos.conf
[Unit]
RequiresMountsFor=/etc/docker /mnt/data /var/lib/docker

[Service]
# env file must be created by docker-prepare.service
EnvironmentFile=/run/dockerd.env
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --containerd=/run/containerd/containerd.sock $DOCKERD_FLAGS
OOMScoreAdjust=-400
