# buildroot-external/rootfs-overlay/usr/lib/systemd/system/systemd-timedated.service.d/localtime.conf
[Service]
Environment=SYSTEMD_ETC_LOCALTIME=/mnt/overlay/etc/localtime
ReadWritePaths=/etc /mnt/overlay/etc
