# buildroot-external/rootfs-overlay/etc/systemd/system/docker.service.d/failure.conf
[Unit]
FailureAction=reboot

[Service]
ExecStopPost=/usr/libexec/docker-failure
