# buildroot-external/rootfs-overlay/etc/systemd/system/docker.service.d/hassos.conf
[Unit]
RequiresMountsFor=/etc/docker /var/lib/docker

[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --storage-driver=overlay2 --log-driver=journald
