# buildroot-external/rootfs-overlay/etc/systemd/system/docker.service.d/hassio.conf
[Unit]
RequiresMountsFor=/etc/docker /var/lib/docker

[Service]
ExecStart=--storage-driver=overlay2 --log-driver=journald
