# buildroot-external/rootfs-overlay/etc/systemd/system/docker.service.d/hassio.conf
[Unit]
RequiresMountsFor=/etc/docker /var/lib/docker
