# buildroot-external/rootfs-overlay/etc/systemd/system/docker.service.d/hassos.conf
[Unit]
RequiresMountsFor=/etc/docker /mnt/data /var/lib/docker

[Service]
OOMScoreAdjust=-400
